# remove_server/kiln_parser/crontab
*/5 * * * * /usr/local/bin/python /app/worker.py 
